# ansible/server_cm_roles/common/files/etc/cron.daily/proxsmtpd_restart
#!/bin/sh
service proxsmtp --full-restart 2>&1 | logger -t PROXSMTPD_RESTART
